# native/controller/config/tftpdd
#!/sbin/openrc-run

command="/usr/sbin/tftpdd"
command_args="-l :1040"
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"

depend() {
    need net
}
